# overlays/usb_gadget/usb_mtp.service
[Unit]
Description=Set up MTP for USB

[Service]
ExecStart=/bin/bash /root/usb_gadget/usb-otg-mtp.sh

[Install]
WantedBy=multi-user.target

